# server/visionserver.init
#!/bin/bash

### BEGIN INIT INFO
# Provides:        visionserver
# Required-Start:  $network
# Required-Stop:   $network
# Should-Start:    
# Default-Start:   2 3 4 5
# Default-Stop:    1
# Short-Description: vision server
# Description: LigerBots vision server
### END INIT INFO

/home/odroid/VisionServer/server/start_server "$1"
